# templates/default/plone.init.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          <%= @name %>
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Should-Start:      <%= @name %>
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start <%= @name %> at boot time
# Description:       Start <%= @name %> at boot time
#
#
#
#
### END INIT INFO

case "$1" in
start)
        echo -n " Starting <%= @name %> "
        <%= @home %>/bin/plonectl start
        ;;
stop)
        echo -n " Stopping <%= @name %> "
        <%= @home %>/bin/plonectl stop
        ;;
restart)
        echo -n " Restarting <%= @name %> "
        <%= @home %>/bin/plonectl restart
        ;;
status)
        echo -n " <%= @name %> Status "
        <%= @home %>/bin/plonectl status
        ;;
*)
        echo "Usage: `basename $0` {start|stop|restart|status}" >&2
        ;;
esac

exit 0
